# supervisord.conf
[supervisord]
nodaemon=true

[program:ollama]
command=ollama serve
autorestart=true

[program:app]
command=python run.py
autorestart=true
startsecs=5
